# scripts/kontrol-acp-crush-adapter.service
[Unit]
Description=Kontrol ACP CRUSH adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Kontrol
EnvironmentFile=-%h/.kontrol/acp-crush-adapter.env
Environment=ACP_ADAPTER_PORT=9877
ExecStart=/usr/bin/env node %h/Kontrol/scripts/acp-crush-adapter.mjs
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
